# misc/systemd/user/org.deepin.dde.Launcher1.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=DDE launcher service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-launcher
Type=dbus
BusName=org.deepin.dde.Launcher1
Restart=on-failure
RestartSec=1s
